# etc/xdg/autostart/regataos-help-select-language.desktop
[Desktop Entry]
Exec=sudo /opt/regataos-help/scripts/select-language
OnlyShowIn=KDE;
Type=Service
MimeType=
Terminal=false
X-KDE-autostart-after=root
X-KDE-StartupNotify=false
X-KDE-autostart-condition=AutoStart:true
Icon=system-run
GenericName=regataos-help-select-language
Name=regataos-help-select-language
